# scripts/launchd/com.engram.ideaspark.feedback.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.engram.ideaspark.feedback</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hani/Box Sync/CLAUDE/engram/scripts/ideaspark_cron.sh</string>
        <string>feedback</string>
    </array>

    <!-- 10:00 PM PT daily -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/hani/Box Sync/CLAUDE/engram/logs/ideaspark_feedback.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hani/Box Sync/CLAUDE/engram/logs/ideaspark_feedback.err.log</string>

</dict>
</plist>
